# agent/data/confd.conf
[supervisord]
user=root
pidfile=/var/run/supervisord.pid
logfile=/var/run/supervisord.log

[program:train]
directory = /usr/src/app
command = python3 ./train_server.py
autostart = true
user=root
stdout_logfile=./runs/train_server_std.log
stderr_logfile=./runs/train_server_err.log
pidfile=/var/run/train_server.pid
numprocs=1

[program:tensorboard]
directory = /usr/src/app
command = python3 ./tensorboard_server.py
autostart = true
user=root
stdout_logfile=./runs/tensorboard_std.log
stderr_logfile=./runs/tensorboard_err.log
pidfile=/var/run/tensorboard.pid
numprocs=1

[program:ui_server]
directory = /usr/src/app
command = python3 ./ui_server.py
autostart = true
user=root
stdout_logfile=./runs/ui_server_std.log
stderr_logfile=./runs/ui_server_err.log
pidfile=/var/run/ui_server.pid
numprocs=1



# supervisord -c confd.conf


